# deploy/systemd/youlya-deploy-watch.service
[Unit]
Description=Youlya pull-based deploy watcher
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/root/youlya
EnvironmentFile=-/root/youlya/.env.production
ExecStart=/bin/bash /root/youlya/scripts/watch-and-deploy.sh
User=root
