# deploy/systemd/youlya-deploy-watch.service
[Unit]
Description=Youlya pull-based deploy watcher
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/youlya
EnvironmentFile=-/opt/youlya/.env.production
ExecStart=/bin/bash /opt/youlya/scripts/watch-and-deploy.sh
User=root
